# smcr_cloud_ha/rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash
exec mysqld_safe --datadir=/data/mariadb --user=mysql 2>&1
